# practice.conf
[program:practice]
directory=/root/test/practice_bot/
command=python3 bot.py
autostart=true
autorestart=true